# src/etc/supervisor/conf.d/confd.conf
[program:confd]
command=/confd -backend etcd -node 127.0.0.1:4001 -config-file /etc/confd/conf.d/logstash.toml -verbose=true
priority=10
numprocs=1
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true