# src/etc/supervisor/conf.d/logstash.conf
[program:logstash]
command=/opt/logstash/bin/logstash -agent -f /etc/logstash/conf.d/logstash.conf
numprocs=1
autostart=true
auto_restart=true
stdout_events_enabled=true
stderr_events_enabled=true
